# deploy/systemd/rawlead-radar.service
[Unit]
Description=RawLead Site radar (main.py + tg_main.py)
After=network-online.target rawlead-api.service
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=all
User=rawlead
Group=rawlead
WorkingDirectory=/opt/rawlead
EnvironmentFile=/opt/rawlead/.env
EnvironmentFile=-/opt/rawlead/.env.site
Environment=RADAR_PROFILE=site
Environment=PYTHONPATH=/opt/rawlead/src
ExecStart=/opt/rawlead/deploy/run-radar-site.sh
Restart=on-failure
RestartSec=30
WatchdogSec=660
MemoryMax=1400M
MemoryHigh=1200M
OOMScoreAdjust=200
KillMode=mixed

[Install]
WantedBy=multi-user.target
